# archived_systems/4runr-agents/scripts/crontab
# Run the pipeline every 6 hours
0 */6 * * * docker exec 4runr-pipeline python /app/run_pipeline.py >> /logs/pipeline.log 2>&1

# Run the scraper every 4 hours
0 */4 * * * docker exec 4runr-scraper python /app/app.py --run-once >> /logs/scraper.log 2>&1

# Run the enricher every 2 hours
0 */2 * * * docker exec 4runr-enricher python /app/app.py --run-once >> /logs/enricher.log 2>&1

# Run the engager every hour
0 * * * * docker exec 4runr-engager python /app/app.py --run-once >> /logs/engager.log 2>&1

# Clean up logs older than 7 days at midnight
0 0 * * * find /logs -type f -name "*.log" -mtime +7 -delete